# debian/scorekeeper/DEBIAN/postinst
#!/bin/sh
/usr/sbin/update-rc.d scorekeeper defaults
echo "addn-hosts=/tmp/scorekeeper.hosts" > /etc/dnsmasq.d/scorekeeper
/etc/init.d/lighttpd restart
/etc/init.d/dnsmasq  restart
