# debian/scorekeeper/etc/network/if-up.d/scorekeeper
#!/bin/sh
/etc/init.d/scorekeeper restart
